# etc/supervisord.d/clean_cache.ini
[program:clean_cache]
command=/opt/opentsdb/bin/clean_cache.sh
stopwaitsecs=60
prioirty=11
autorestart=true
