# etc/supervisord.d/feed_stats.ini
[program:feed_stats]
command=/opt/opentsdb/bin/feed_stats.sh
stopwaitsecs=60
prioirty=11
autorestart=true
