# etc/supervisord.d/tsdb.ini
[program:opentsdb]
command=/opt/opentsdb/bin/start_tsdb.sh
stopwaitsecs=60
prioirty=10
autorestart=true
